# yanic/files/yanic-to-zonefile-reload.path
[Path]
PathChanged=/var/lib/knot/zones/nodes.ffrn.de.zone

[Install]
WantedBy=multi-user.target
