# yanic/files/yanic-to-zonefile-reload.path
[Path]
PathExists=/var/lib/knot/nodes.ffrn.de.zone
PathChanged=/var/lib/knot/nodes.ffrn.de.zone

[Install]
WantedBy=multi-user.target
